# .crontab
 */5 * * * * env DISPLAY=:0 /home/nick/.local/bin/_batterycheck
