# .crontab
 */5 * * * * env DISPLAY=:0 _batterycheck
 */20 * * * * crontab -l > $HOME/.crontab && git add $HOME/.crontab
