# .crontab
*/5 * * * * env DISPLAY=:0 `which _batterycheck`
