# .crontab
 */5 * * * * env DISPLAY=:0 `which _batterycheck`
 */20 * * * * crontab -l > $HOME/.crontab
