# workspace/systemd/system/batterycheck.service
[Unit]
Description=battery checker

[Service]
User=nick
Type=simple
ExecStart=/usr/bin/_batterycheck
