# workspace/systemd/system/batterycheck.service
[Unit]
Description=battery checker

[Service]
Type=simple
ExecStart=/usr/bin/_batterycheck

