# workspace/systemd/system/batterycheck.timer
[Unit]
Description=battery checker

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
