# workspace/systemd/system/batterycheck.timer
[Unit]
Description=battery checker

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
